# stravainkydashboard.timer
[Timer]
OnCalendar=*-*-* *:00:00
Unit=stravainkydashboard.service
